# api/flockyou.service
[Unit]
Description=Flock-You detection server
After=network.target

[Service]
Type=simple
# User= must match the actual Pi username. pi-setup.sh patches this automatically.
User=pi
WorkingDirectory=/home/pi/flock-you/api
ExecStart=/usr/bin/python3 flockyou.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flockyou

[Install]
WantedBy=multi-user.target
