# api/flockyou.service
[Unit]
Description=Flock-You detection server
After=network.target

[Service]
Type=simple
User=pi  # TODO: change to match your Pi username before running pi-setup.sh
WorkingDirectory=/home/pi/flock-you/api
ExecStart=/usr/bin/python3 flockyou.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flockyou

[Install]
WantedBy=multi-user.target
